# gentoo/init.d/pidiod.in
#!/sbin/runscript
# -*- shell-script -*-
# Copyright 2010 Bent Bisballe Nyeng (deva@aasimon.org)
# Distributed under the terms of the GNU General Public License v2

_PROG="Pidiod"
_PREFIX="@prefix@"
_CONF="${_PREFIX}/etc/pidiod.conf"
_BIN="${_PREFIX}/bin/pidiod"
_LOG="/var/log/pidiod.log"
_PID="/var/run/pidiod.pid"

depend() {
    need net
    need postgresql
#    provide pidiod
}

checkconfig() {
    if [ ! -f ${_CONF} ] ; then
	eerror "Please create ${_CONF}. See ${_CONF}.example for examples."
	return 1
    fi
    return 0
}

start() {
    checkconfig || return $?
    
    ebegin "Starting ${_PROG}"
    start-stop-daemon --start --exec ${_BIN} \
	--pidfile ${_PID} \
	-- --pidfile ${_PID} --config ${_CONF} \
	-D +fixme,+err,+warn,+info --logfile ${_LOG}
    eend $? "Failed to start ${_PROG}"
}

stop() {
    ebegin "Stopping ${_PROG}"
    start-stop-daemon --stop -s INT \
	--pidfile ${_PID} \
	--exec ${_BIN}
    eend $? "Failed to stop ${_PROG}"
}
